# crontab
0 0 * * * /usr/local/bin/python3 /code/project/etl.py > /proc/1/fd/1 2>/proc/1/fd/2
